# .config/systemd/user/polkit.service
[Unit]
Description=Polkit authentication agents

[Service]
ExecStart=/usr/lib/polkit-1/polkit-agent-helper-1

[Install]
WantedBy=multi-user.target
